# etc/snap-yearly.timer
[Unit]
Description=Trigger snap yearly snapshots

[Timer]
OnCalendar=yearly
Persistent=true

[Install]
WantedBy=multi-user.target
