# etc/snap-monthly.timer
[Unit]
Description=Trigger snap monthly snapshots

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=multi-user.target
