# etc/snap-daily.timer
[Unit]
Description=Trigger snap daily snapshots

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
