# etc/snap-service-template.service
[Unit]
Description=Daily btrfs snapshot

[Service]
Type=oneshot
ExecStart=/usr/bin/snap-create -t daily /home /home/.snap
