# etc/snap-prune.timer
[Unit]
Description=Trigger for snap-prune

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=multi-user.target
