# etc/snap-prune-service-template.service
[Unit]
Description=Prune btrfs snapshots

[Service]
Type=oneshot
ExecStart=/usr/bin/snap-prune -H 12 -d 14 -w 4 /home/.snap
