# etc/snap-hourly.timer
[Unit]
Description=Trigger snap hourly snapshots

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=multi-user.target
