# etc/snap-weekly.timer
[Unit]
Description=Trigger snap weekly snapshots

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=multi-user.target
